# dozzle/rootfs/etc/services.d/dozzle/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App — Dozzle
# Ingress : --base depuis bashio::addon.ingress_entry
# Agent optionnel : expose cet hôte aux autres instances Dozzle (dozzle agent)
# s6-overlay : https://github.com/just-containers/s6-overlay
# ==============================================================================

set -e

declare -a DOZZLE_ARGS
declare DOZZLE_BASE_PATH
declare DOZZLE_ADDR
declare LOG_LEVEL
declare FILTER
declare NO_ANALYTICS
declare ENABLE_ACTIONS
declare INGRESS_ENTRY
declare REMOTE_AGENTS
declare AGENT_PORT
declare ENABLE_AGENT

# --- Docker socket (docker_api: true dans config.yaml) ---
if [ -S /var/run/docker.sock ]; then
    bashio::log.info "Docker socket présent"
    if curl -s --unix-socket /var/run/docker.sock http://localhost/version >/dev/null 2>&1; then
        DOCKER_VERSION="$(curl -s --unix-socket /var/run/docker.sock http://localhost/version | jq -r '.Version')"
        bashio::log.info "API Docker OK — version ${DOCKER_VERSION}"
    else
        bashio::log.warning "Socket Docker présent mais l’API ne répond pas"
    fi
else
    bashio::log.error "Pas de socket à /var/run/docker.sock — activez docker_api dans le manifest"
fi

# --- Ingress (proxy Supervisor / intégration panel frontend) ---
INGRESS_ENTRY="$(bashio::addon.ingress_entry)"
bashio::log.info "Ingress entry: ${INGRESS_ENTRY}"
DOZZLE_BASE_PATH="${INGRESS_ENTRY}"
bashio::log.info "Dozzle --base: ${DOZZLE_BASE_PATH}"

DOZZLE_ADDR=":8080"

LOG_LEVEL="$(bashio::config 'log_level' 'info')"
FILTER="$(bashio::config 'filter' '')"
NO_ANALYTICS="$(bashio::config 'no_analytics' 'true')"
ENABLE_ACTIONS="$(bashio::config 'enable_actions' 'false')"
ENABLE_AGENT="$(bashio::config 'enable_agent' 'false')"
AGENT_PORT="$(bashio::config 'agent_port' '7007')"
REMOTE_AGENTS="$(bashio::config 'remote_agents' '')"

DOZZLE_ARGS=(
    "--addr" "${DOZZLE_ADDR}"
    "--base" "${DOZZLE_BASE_PATH}"
    "--level" "${LOG_LEVEL}"
)

if bashio::var.true "${NO_ANALYTICS}"; then
    DOZZLE_ARGS+=("--no-analytics")
fi

if bashio::var.has_value "${FILTER}"; then
    DOZZLE_ARGS+=("--filter" "${FILTER}")
    bashio::log.info "Filtre conteneurs: ${FILTER}"
fi

if bashio::var.true "${ENABLE_ACTIONS}"; then
    DOZZLE_ARGS+=("--enable-actions")
    bashio::log.info "Actions conteneur (restart/stop) activées"
fi

DOZZLE_ARGS+=("--auth-provider" "none")

if bashio::var.has_value "${REMOTE_AGENTS}"; then
    export DOZZLE_REMOTE_AGENT="${REMOTE_AGENTS}"
    bashio::log.info "Agents distants (DOZZLE_REMOTE_AGENT): ${REMOTE_AGENTS}"
fi

export DOZZLE_BASE="${DOZZLE_BASE_PATH}"
export DOZZLE_ADDR="${DOZZLE_ADDR}"
export DOZZLE_LEVEL="${LOG_LEVEL}"
export DOZZLE_AUTH_PROVIDER="none"
if bashio::var.true "${NO_ANALYTICS}"; then
    export DOZZLE_NO_ANALYTICS="true"
fi

if bashio::var.true "${ENABLE_AGENT}"; then
    bashio::log.info "Démarrage du Dozzle Agent sur 0.0.0.0:${AGENT_PORT}..."
    if bashio::var.has_value "$(bashio::config 'agent_hostname')"; then
        AH="$(bashio::config 'agent_hostname')"
        bashio::log.info "Nom d’hôte agent: ${AH}"
        env DOZZLE_HOSTNAME="${AH}" /usr/local/bin/dozzle agent --addr "0.0.0.0:${AGENT_PORT}" --level "${LOG_LEVEL}" &
    else
        /usr/local/bin/dozzle agent --addr "0.0.0.0:${AGENT_PORT}" --level "${LOG_LEVEL}" &
    fi
    sleep 1
fi

bashio::log.info "============================================"
bashio::log.info " Dozzle — UI logs Docker"
bashio::log.info " Base (Ingress): ${DOZZLE_BASE_PATH}"
bashio::log.info " Écoute         : ${DOZZLE_ADDR}"
bashio::log.info " Niveau log     : ${LOG_LEVEL}"
bashio::log.info " Filtre         : ${FILTER:-aucun}"
bashio::log.info " Agent intégré  : ${ENABLE_AGENT} (port ${AGENT_PORT})"
bashio::log.info "============================================"

bashio::log.info "Commande: /usr/local/bin/dozzle ${DOZZLE_ARGS[*]}"
exec /usr/local/bin/dozzle "${DOZZLE_ARGS[@]}"
