# dozzle/rootfs/etc/services.d/dozzle/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App — Dozzle
# Ingress  : --base from bashio::addon.ingress_entry
# Agent    : optional built-in dozzle agent (enable_agent option)
# s6-overlay: https://github.com/just-containers/s6-overlay
# ==============================================================================

set -e

declare -a DOZZLE_ARGS
declare DOZZLE_BASE_PATH
declare DOZZLE_ADDR
declare LOG_LEVEL
declare FILTER
declare NO_ANALYTICS
declare ENABLE_ACTIONS
declare INGRESS_ENTRY
declare REMOTE_AGENTS
declare AGENT_PORT
declare AGENT_HOSTNAME
declare ENABLE_AGENT

# --- Docker socket (docker_api: true in config.yaml) -------------------------
if [ -S /var/run/docker.sock ]; then
    if curl -s --unix-socket /var/run/docker.sock http://localhost/version >/dev/null 2>&1; then
        DOCKER_VERSION="$(curl -s --unix-socket /var/run/docker.sock http://localhost/version | jq -r '.Version')"
        bashio::log.info "Docker API OK — engine version ${DOCKER_VERSION}"
    else
        bashio::log.warning "Docker socket found but API is not responding"
    fi
else
    bashio::log.error "No socket at /var/run/docker.sock — enable docker_api in the manifest"
fi

# --- Ingress (Supervisor proxy / panel frontend) ------------------------------
# The Supervisor STRIPS the /api/hassio_ingress/TOKEN prefix before forwarding
# to the container — the container always receives requests at "/" (ingress_entry: /).
# Do NOT pass the full ingress token as --base; use "/" (Dozzle default).
INGRESS_ENTRY="$(bashio::addon.ingress_entry)"
DOZZLE_BASE_PATH="/"
bashio::log.info "Ingress entry (external): ${INGRESS_ENTRY}"
bashio::log.info "Dozzle base path (container): ${DOZZLE_BASE_PATH}"

DOZZLE_ADDR=":8081"  # nginx listens on :8080 and proxies here

# --- Read options from HA config ---------------------------------------------
LOG_LEVEL="$(bashio::config 'log_level' 'info')"
FILTER="$(bashio::config 'filter' '')"
NO_ANALYTICS="$(bashio::config 'no_analytics' 'true')"
ENABLE_ACTIONS="$(bashio::config 'enable_actions' 'false')"
ENABLE_AGENT="$(bashio::config 'enable_agent' 'false')"
AGENT_PORT="$(bashio::config 'agent_port' '7007')"
AGENT_HOSTNAME="$(bashio::config 'agent_hostname' '')"
REMOTE_AGENTS="$(bashio::config 'remote_agents' '')"

# --- Build Dozzle argument list ----------------------------------------------
DOZZLE_ARGS=(
    "--addr"  "${DOZZLE_ADDR}"
    "--base"  "${DOZZLE_BASE_PATH}"
    "--level" "${LOG_LEVEL}"
    "--auth-provider" "none"
)

if bashio::var.true "${NO_ANALYTICS}"; then
    DOZZLE_ARGS+=("--no-analytics")
fi

if bashio::var.has_value "${FILTER}"; then
    DOZZLE_ARGS+=("--filter" "${FILTER}")
    bashio::log.info "Container filter: ${FILTER}"
fi

if bashio::var.true "${ENABLE_ACTIONS}"; then
    DOZZLE_ARGS+=("--enable-actions")
    bashio::log.info "Container actions (restart/stop) enabled"
fi

# --- Export env vars used by Dozzle ------------------------------------------
export DOZZLE_BASE="${DOZZLE_BASE_PATH}"
export DOZZLE_ADDR="${DOZZLE_ADDR}"  # nginx :8080 → Dozzle :8081
export DOZZLE_LEVEL="${LOG_LEVEL}"
export DOZZLE_AUTH_PROVIDER="none"

if bashio::var.true "${NO_ANALYTICS}"; then
    export DOZZLE_NO_ANALYTICS="true"
fi

# --- Remote agents -----------------------------------------------------------
if bashio::var.has_value "${REMOTE_AGENTS}"; then
    export DOZZLE_REMOTE_AGENT="${REMOTE_AGENTS}"
    bashio::log.info "Remote agents (DOZZLE_REMOTE_AGENT): ${REMOTE_AGENTS}"
fi

# --- Built-in agent (optional) -----------------------------------------------
if bashio::var.true "${ENABLE_AGENT}"; then
    bashio::log.info "Starting built-in Dozzle agent on 0.0.0.0:${AGENT_PORT} ..."
    if bashio::var.has_value "${AGENT_HOSTNAME}"; then
        bashio::log.info "Agent hostname: ${AGENT_HOSTNAME}"
        DOZZLE_HOSTNAME="${AGENT_HOSTNAME}" \
            /usr/local/bin/dozzle agent \
                --addr "0.0.0.0:${AGENT_PORT}" \
                --level "${LOG_LEVEL}" &
    else
        /usr/local/bin/dozzle agent \
            --addr "0.0.0.0:${AGENT_PORT}" \
            --level "${LOG_LEVEL}" &
    fi
    sleep 1
    bashio::log.info "Built-in agent started on port ${AGENT_PORT}"
fi

# --- Summary -----------------------------------------------------------------
bashio::log.info "========================================"
bashio::log.info " Dozzle — Docker log viewer"
bashio::log.info " Ingress base : ${DOZZLE_BASE_PATH}"
bashio::log.info " Listen addr  : ${DOZZLE_ADDR}"
bashio::log.info " Log level    : ${LOG_LEVEL}"
bashio::log.info " Filter       : ${FILTER:-none}"
bashio::log.info " Agent        : ${ENABLE_AGENT} (port ${AGENT_PORT})"
bashio::log.info " Remote agents: ${REMOTE_AGENTS:-none}"
bashio::log.info "========================================"

bashio::log.info "Command: /usr/local/bin/dozzle ${DOZZLE_ARGS[*]}"
exec /usr/local/bin/dozzle "${DOZZLE_ARGS[@]}"
