# dozzle/rootfs/etc/services.d/dozzle/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant App - Dozzle
# Ingress: nginx on :8080 adds back /api/hassio_ingress/TOKEN prefix,
#          Dozzle on :8081 runs with --base TOKEN_PATH so all asset/API
#          URLs are correctly scoped to the Ingress URL.
# Agent:   optional built-in dozzle agent on port 7007 (enable_agent option)
# ==============================================================================

set -e

declare -a DOZZLE_ARGS
declare DOZZLE_ADDR
declare LOG_LEVEL
declare FILTER
declare NO_ANALYTICS
declare ENABLE_ACTIONS
declare INGRESS_ENTRY
declare REMOTE_AGENTS
declare ENABLE_AGENT
declare AGENT_HOSTNAME
declare ENABLE_DIRECT_ACCESS

# --- Docker socket (docker_api: true in config.yaml) -------------------------
if [ -S /var/run/docker.sock ]; then
    if curl -s --unix-socket /var/run/docker.sock http://localhost/version >/dev/null 2>&1; then
        DOCKER_VERSION="$(curl -s --unix-socket /var/run/docker.sock http://localhost/version | jq -r '.Version')"
        bashio::log.info "Docker API OK - engine ${DOCKER_VERSION}"
    else
        bashio::log.warning "Docker socket found but API is not responding"
    fi
else
    bashio::log.error "No socket at /var/run/docker.sock - enable docker_api in the manifest"
fi

# --- Ingress base path -------------------------------------------------------
# Supervisor STRIPS /api/hassio_ingress/TOKEN before forwarding to the container.
# We pass the full token path as --base so Dozzle embeds it in all asset/API URLs.
# nginx adds the prefix back (dozzle.conf rewrite rule).
INGRESS_ENTRY="$(bashio::addon.ingress_entry)"

# nginx on :8080 → Dozzle on :8081 (internal only)
DOZZLE_ADDR=":8081"

# --- Read options from HA config ---------------------------------------------
LOG_LEVEL="$(bashio::config 'log_level' 'info')"
FILTER="$(bashio::config 'filter' '')"
NO_ANALYTICS="$(bashio::config 'no_analytics' 'true')"
ENABLE_ACTIONS="$(bashio::config 'enable_actions' 'false')"
ENABLE_AGENT="$(bashio::config 'enable_agent' 'false')"
AGENT_HOSTNAME="$(bashio::config 'agent_hostname' '')"
REMOTE_AGENTS="$(bashio::config 'remote_agents' '')"
ENABLE_DIRECT_ACCESS="$(bashio::config 'enable_direct_access' 'false')"

# --- Build Dozzle argument list ----------------------------------------------
DOZZLE_ARGS=(
    "--addr"          "${DOZZLE_ADDR}"
    "--base"          "${INGRESS_ENTRY}"
    "--level"         "${LOG_LEVEL}"
    "--auth-provider" "none"
)

if bashio::var.true "${NO_ANALYTICS}"; then
    DOZZLE_ARGS+=("--no-analytics")
fi

if bashio::var.has_value "${FILTER}"; then
    DOZZLE_ARGS+=("--filter" "${FILTER}")
fi

if bashio::var.true "${ENABLE_ACTIONS}"; then
    DOZZLE_ARGS+=("--enable-actions")
fi

# --- Export env vars used by Dozzle ------------------------------------------
export DOZZLE_BASE="${INGRESS_ENTRY}"
export DOZZLE_ADDR="${DOZZLE_ADDR}"
export DOZZLE_LEVEL="${LOG_LEVEL}"
export DOZZLE_AUTH_PROVIDER="none"

if bashio::var.true "${NO_ANALYTICS}"; then
    export DOZZLE_NO_ANALYTICS="true"
fi

if bashio::var.has_value "${REMOTE_AGENTS}"; then
    export DOZZLE_REMOTE_AGENT="${REMOTE_AGENTS}"
fi

# --- Direct access instance (optional, port 8082 → nginx 8088) ---------------
# Runs a second Dozzle with --base / so browsers can access it directly
# without the ingress token prefix. nginx serves it on port 8088.
# Security: only activated when enable_direct_access is true in config.
if bashio::var.true "${ENABLE_DIRECT_ACCESS}"; then
    bashio::log.info "Starting direct-access Dozzle on :8082 (nginx will serve on :8088)"
    declare -a DIRECT_ARGS=(
        "--addr"          ":8082"
        "--base"          "/"
        "--level"         "${LOG_LEVEL}"
        "--auth-provider" "none"
    )
    if bashio::var.true "${NO_ANALYTICS}"; then
        DIRECT_ARGS+=("--no-analytics")
    fi
    if bashio::var.has_value "${FILTER}"; then
        DIRECT_ARGS+=("--filter" "${FILTER}")
    fi
    if bashio::var.true "${ENABLE_ACTIONS}"; then
        DIRECT_ARGS+=("--enable-actions")
    fi
    if bashio::var.has_value "${REMOTE_AGENTS}"; then
        DOZZLE_REMOTE_AGENT="${REMOTE_AGENTS}" \
            /usr/local/bin/dozzle "${DIRECT_ARGS[@]}" &
    else
        /usr/local/bin/dozzle "${DIRECT_ARGS[@]}" &
    fi
fi

# --- Built-in agent (optional, port 7007) ------------------------------------
# This exposes HA's Docker to OTHER Dozzle instances on your network.
# Map port 7007 in the Network tab to make it reachable from outside.
# You do NOT need this to see HA containers in this add-on.
if bashio::var.true "${ENABLE_AGENT}"; then
    bashio::log.info "Starting built-in agent on 0.0.0.0:7007 (map port 7007 in Network tab)"
    if bashio::var.has_value "${AGENT_HOSTNAME}"; then
        bashio::log.info "Agent display name: ${AGENT_HOSTNAME}"
        DOZZLE_HOSTNAME="${AGENT_HOSTNAME}" \
            /usr/local/bin/dozzle agent --addr "0.0.0.0:7007" --level "${LOG_LEVEL}" &
    else
        /usr/local/bin/dozzle agent --addr "0.0.0.0:7007" --level "${LOG_LEVEL}" &
    fi
    sleep 1
    bashio::log.info "Built-in agent ready on port 7007"
fi

# --- Startup banner ----------------------------------------------------------
bashio::log.info "================================================"
bashio::log.info "  Dozzle for Home Assistant"
bashio::log.info "  App version  : ${BUILD_VERSION:-unknown}"
bashio::log.info "  Dozzle binary: ${HA_DOZZLE_BIN_VERSION:-unknown}"
bashio::log.info "------------------------------------------------"
bashio::log.info "  Ingress base : ${INGRESS_ENTRY}"
bashio::log.info "  Listen addr  : ${DOZZLE_ADDR} (nginx :8080 → here)"
bashio::log.info "  Log level    : ${LOG_LEVEL}"
bashio::log.info "  Filter       : ${FILTER:-none}"
bashio::log.info "  Actions      : ${ENABLE_ACTIONS}"
bashio::log.info "  Agent        : ${ENABLE_AGENT} (port 7007, name: ${AGENT_HOSTNAME:-default})"
bashio::log.info "  Remote agents: ${REMOTE_AGENTS:-none}"
bashio::log.info "  Direct access: ${ENABLE_DIRECT_ACCESS} (port 8088, map in Network tab)"
bashio::log.info "================================================"

exec /usr/local/bin/dozzle "${DOZZLE_ARGS[@]}"
