# dozzle/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# nginx - reverse proxy in front of Dozzle
# Listens on :8080 (Ingress-facing), forwards to Dozzle on :8081.
#
# Key: Supervisor STRIPS /api/hassio_ingress/TOKEN before forwarding.
#      Dozzle runs with --base TOKEN_PATH so asset/API URLs in the HTML embed
#      the full token path. nginx adds the prefix back so Dozzle sees it.
# ==============================================================================

INGRESS_ENTRY="$(bashio::addon.ingress_entry)"
ENABLE_DIRECT_ACCESS="$(bashio::config 'enable_direct_access' 'false')"

# Inject the actual ingress base path into the nginx config template
sed -i "s|__INGRESS_BASE__|${INGRESS_ENTRY}|g" /etc/nginx/http.d/dozzle.conf

# Write direct access server block to a separate conf file when enabled.
# nginx loads all *.conf from http.d, so direct_access.conf is picked up automatically.
if bashio::var.true "${ENABLE_DIRECT_ACCESS}"; then
    cat > /etc/nginx/http.d/direct_access.conf << 'NGINX_EOF'
# Direct access - generated by /etc/services.d/nginx/run
# Proxies to a second Dozzle instance running with --base / on port 8082.
server {
    listen 8088;
    server_name _;
    access_log off;

    location / {
        proxy_pass         http://127.0.0.1:8082;
        proxy_http_version 1.1;
        proxy_set_header   Host $host;
        proxy_set_header   Connection '';
        proxy_buffering    off;
        proxy_cache        off;
        proxy_read_timeout 86400s;
    }
}
NGINX_EOF
    bashio::log.info "nginx: direct access enabled on :8088 → Dozzle :8082"
else
    rm -f /etc/nginx/http.d/direct_access.conf
fi

mkdir -p /tmp/nginx_client /tmp/nginx_proxy

bashio::log.info "nginx: ingress base = ${INGRESS_ENTRY}"
bashio::log.info "nginx: :8080  rewrite → Dozzle :8081${INGRESS_ENTRY}"

# Wait for Dozzle to accept connections before nginx starts
bashio::log.info "nginx: waiting for Dozzle on :8081..."
for i in $(seq 1 30); do
    curl -s --max-time 1 -o /dev/null http://127.0.0.1:8081/ 2>/dev/null && break || true
    sleep 1
done
bashio::log.info "nginx: Dozzle ready - starting proxy"

exec nginx -g 'daemon off;' -e /dev/stderr
